# deploy/crontab
SHELL=/bin/bash

00 00 * * * cd /var/scpdb/featured/current && /usr/bin/node build rotateDaily
02 00 * * 1 cd /var/scpdb/featured/current && /usr/bin/node build rotateWeekly
